# inittab
::respawn:/sbin/syslogd -n -O -
::respawn:/usr/sbin/crond -f -l 9
::respawn:/usr/sbin/named -f -4 -u named -t /opt/named
