# inittab
::wait:/usr/local/bin/startup
::respawn:/usr/local/bin/start_syslogd
::respawn:/usr/sbin/crond -f -l 9
::respawn:/usr/local/bin/start_bind
